# ansible/roles/museumplus_pas_workflow/templates/config/stop-sync-attachments.timer.j2
[Unit]
Description=Start stop-sync-attachments

[Timer]
OnCalendar={{ sync_attachments_end_date }}

[Install]
WantedBy=timers.target
