# ansible/roles/museumplus_pas_workflow/templates/config/sync-processed-sips.service.j2
[Unit]
Description=Synchronize processed SIPS from the DPRES service
After=network.target

[Service]
Type=simple
WorkingDirectory={{ museumplus_pas_workflow_install_path }}
Environment=LANG=en_US.UTF-8
ExecStart={{ museumplus_pas_workflow_install_path }}/run-venv.sh sync-processed-sips
Restart=on-failure
User={{ museumplus_pas_user }}
Group={{ museumplus_pas_user }}

[Install]
WantedBy=multi-user.target
