# ansible/roles/museumplus_pas_workflow/templates/config/sync-hashes.timer.j2
[Unit]
Description=Start sync-hashes

[Timer]
OnCalendar={{ sync_hashes_start_date }}

[Install]
WantedBy=timers.target
