# ansible/roles/museumplus_pas_workflow/templates/config/sync-processed-sips.timer.j2
[Unit]
Description=Start sync-processed-sips

[Timer]
OnCalendar={{ sync_processed_sips_start_date }}

[Install]
WantedBy=timers.target
