# ansible/roles/museumplus_pas_backup/templates/backup-postgresql.timer.j2
[Unit]
Description=Perform a daily MuseumPlus PAS PostgreSQL backup

[Timer]
OnCalendar={{ museumplus_pas_backup_start_date }}

[Install]
WantedBy=timers.target
