# ansible/roles/passari_backup/templates/backup-postgresql.timer.j2
[Unit]
Description=Perform a daily Passari PostgreSQL backup

[Timer]
OnCalendar={{ passari_backup_start_date }}

[Install]
WantedBy=timers.target
