# ansible/roles/passari_workflow/templates/config/passari-worker-NAME@.timer.j2
[Unit]
Description=Timer to start Passari %j worker %i

[Timer]
OnCalendar={{download_objects_start_date}}

[Install]
WantedBy=timers.target
