# ansible/roles/passari_workflow/templates/config/passari-worker-stop-NAME@.service.j2
[Unit]
Description=Timer to stop Passari %j worker %i

[Service]
ExecStart=/usr/bin/systemctl stop passari-worker-%j@%i.service
