# ansible/roles/passari_workflow/templates/config/passari-worker-stop-NAME@.timer.j2
[Unit]
Description=Stop Passari %j worker %i

[Timer]
OnCalendar={{download_objects_end_date}}

[Install]
WantedBy=timers.target
