# ansible/roles/passari_workflow/templates/config/stop-sync-objects.timer.j2
[Unit]
Description=Start stop-sync-objects

[Timer]
OnCalendar={{ sync_objects_end_date }}

[Install]
WantedBy=timers.target
