# ansible/roles/passari_workflow/templates/config/stop-download-object-worker@.timer.j2
[Unit]
Description=Stop download-object worker %i

[Timer]
OnCalendar={{ download_objects_end_date }}

[Install]
WantedBy=timers.target
