# ansible/roles/passari_workflow/templates/config/stop-sync-objects.service.j2
[Unit]
Description=Stop sync-objects (timer activated)
After=network.target

[Service]
ExecStart=/usr/bin/systemctl stop sync-objects.service

[Install]
WantedBy=multi-user.target
