# ansible/roles/passari_workflow/templates/config/sync-attachments.timer.j2
[Unit]
Description=Start sync-attachments

[Timer]
OnCalendar={{ sync_attachments_start_date }}

[Install]
WantedBy=timers.target
