# ansible/roles/passari_workflow/templates/config/sync-objects.timer.j2
[Unit]
Description=Start sync-objects

[Timer]
OnCalendar={{sync_objects_start_date}}

[Install]
WantedBy=timers.target
