# ansible/roles/passari_workflow/templates/config/sync-processed-sips.service.j2
[Unit]
Description=Synchronize processed SIPS from the DPRES service
After=network.target

[Service]
Type=simple
WorkingDirectory={{passari_workflow_install_path}}
Environment=LANG=en_US.UTF-8
ExecStart={{passari_workflow_install_path}}/run-venv.sh sync-processed-sips
Restart=on-failure
User={{passari_user}}
Group={{passari_user}}

[Install]
WantedBy=multi-user.target
